# templates/systemd/openclaw@.service
[Unit]
Description=OpenClaw Agent — %i
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=__REPO_DIR__/agents/%i/.env
Environment=OPENCLAW_STATE_DIR=__SERVER_DATA_DIR__/agents/%i/config
Environment=PATH=__NODE_BIN_DIR__:/usr/local/bin:/usr/bin:/bin
ExecStart=__OPENCLAW_BIN__ gateway
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__SERVER_DATA_DIR__ __REPO_DIR__/agents
PrivateTmp=true

[Install]
WantedBy=default.target
